# supervisord.conf
[supervisord]
nodaemon=true

[program:couchpotato]
command=/opt/couchpotato/CouchPotato.py --console_log --data_dir=/config --config_file=/config/settings.conf